# src/nginx/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /src/app/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx